# supervisord.conf
[supervisord]
nodaemon=true

[program:dcem]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
autorestart=true
command=bash /opt/dcem-supervisord-wrapper.sh
